# deploy/rootcause-cert-renew.service
[Unit]
Description=Renew Let's Encrypt certificates for RootCause
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
ExecStart=/usr/local/bin/rootcause-cert-renew.sh
